# systemd/hermes-release-radar.service
[Unit]
Description=Hermes Release Radar local helper
Documentation=file:%h/.hermes/release-radar/index.html
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.hermes/release-radar
ExecStart=/usr/bin/python3 %h/.hermes/release-radar/serve.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
